# asgard-projects-distribution/src/standard/scripts/api.service
[Unit]
Description=asgard-projects

[Service]
WorkingDirectory=/opt/api/
ExecStart=/opt/api/scripts/api.sh

[Install]
WantedBy=multi-user.target
